# collector/munincollector/samples/etc/cron.d/munin-collector-pickle-plugin-configs
#
# cron-job to reset collectors time ranges.
#

MAILTO=

* * * * *     root test -x /opt/PyramidContainers/VirtualEnvironment-Python-2.6/Munin-Collector/munincollector/utilities/munin-collector-pickle-plugin-configs && /opt/PyramidContainers/VirtualEnvironment-Python-2.6/Munin-Collector/munincollector/utilities/munin-collector-pickle-plugin-configs
